# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user start-up script.
#

WIRELESS_INTERFACE="wlan0"

/sbin/ifconfig ${WIRELESS_INTERFACE} up

/usr/sbin/wpa_supplicant -BDwext -i ${WIRELESS_INTERFACE} -c /etc/wpa_supplicant/wpa_supplicant.conf > /dev/null 2>&1

/sbin/dhcpcd ${WIRELESS_INTERFACE}

#Apple buletooth keyboard
/usr/bin/hciconfig x up && \
	/usr/bin/sleep 1.3 && \
	/usr/bin/hciconfig x reset
